# resources/configs/1-ultroid.conf
[program:ultroid]
command=bash -c '/root/TeamUltroid/startup'
stderr_logfile=/app/.vubuntu/assets/logs/ultroid.stderr.log
priority=1
